# infra/systemd/parrot-cold-backup@.service
[Unit]
Description=Parrot %i stopped cold backup
Documentation=file:/opt/parrot/current/docs/operations/backup-restore-runbook.md
RequiresMountsFor=/srv/project-conversation/%i /mnt/bigboi/project-conversation/%i/backups
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
User=parrot
Group=parrot
UMask=0077
EnvironmentFile=/srv/project-conversation/%i/config/deployment.env
ExecStart=/opt/parrot/current/infra/scripts/backup.sh --env-file /srv/project-conversation/%i/config/deployment.env --apply --confirm project-conversation-%i
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/srv/project-conversation/%i /mnt/bigboi/project-conversation/%i/backups
RestrictSUIDSGID=true
LockPersonality=true
